# target/linux/spacemit/base-files/etc/hotplug.d/iface/99-wwan-promisc
# brcmfmac firmware ARP offload eats ARP frames relayd must answer;
# promiscuous mode disables the firmware filtering. Required for the
# relayd wireless bridge to be loss-free.
[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wwan" ] && ip link set "$DEVICE" promisc on
