# target/linux/spacemit/base-files/etc/uci-defaults/99-rv2-lan
#!/bin/sh
uci -q batch <<UCI
set network.lan.proto=static
set network.lan.ipaddr=192.168.1.250
set network.lan.netmask=255.255.255.0
delete network.lan.dns
add_list network.lan.dns=192.168.1.1
set dhcp.lan.ignore=1
UCI
uci commit network
uci commit dhcp
exit 0
